# lib/systemd/system/deploy.service
[Unit]
Description=Deploy service

[Service]
Type=oneshot
ExecStart=/var/lib/ybs32/bin/deploy
SyslogIdentifier=deploy

[Install]
WantedBy=multi-user.target
