# lib/systemd/system/xserver.service
[Unit]
Description=Display service

[Service]
Environment="HOME=/home/chrome"
ExecStart=/usr/bin/startx
